# deploy/kingo-daemon.service
[Unit]
Description=Kingo Linux VPN Daemon
After=network-online.target
Wants=network-online.target

[Service]
# FIX BUG-025 & 026: Type=notify, Watchdog, and Restart=always
Type=notify
WatchdogSec=30s

User=kingo-vpn
Group=kingo-vpn

ExecStart=/usr/local/bin/kingo-linux-vpn-daemon
Restart=always
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=60

# Security Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/kingo-vpn /run/kingo-vpn /var/log/kingo-vpn
RuntimeDirectory=kingo-vpn
StateDirectory=kingo-vpn
LogsDirectory=kingo-vpn
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
