# packaging/kingo-linux-vpn-tray.service
[Unit]
Description=Kingo Linux VPN Tray Icon
After=graphical-session.target
Wants=graphical-session.target
PartOf=kingo-linux-vpn-daemon.service

[Service]
Type=simple
ExecStart=%h/.local/bin/kingo-tray
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
